# __debian/prerm
#!/usr/bin/env bash

systemctl stop prometheusSDlistener
systemctl disable prometheusSDlistener

