# __alpine/prometheusListener.pre-deinstall
#!/usr/bin/env ash

rc-service prometheusSDlistener stop
rc-update del prometheusSDlistener

exit 0
