# __alpine/prometheusListener.post-install
#!/usr/bin/env ash

chmod 775 /opt/sbin/prometheusSDlistener
mkdir -p /etc/prometheusSDlistener
touch /etc/prometheusSDlistener/prometheusSDlistener.json
chown -R prometheus:prometheus /etc/prometheusSDlistener/prometheusSDlistener.json

rc-update add prometheusSDlistener
rc-service prometheusSDlistener start

exit 0
